# resources/etc/cron.d/freshclam
0 4 * * * /usr/bin/freshclam --quiet --config-file=/etc/clamav/freshclam.conf